# deploy/systemd/fund.service
[Unit]
Description=AI Trading Fund (Claude Agent SDK)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fund
WorkingDirectory=/opt/ai-trading-fund
EnvironmentFile=/opt/ai-trading-fund/.env
ExecStart=/opt/ai-trading-fund/.venv/bin/python -m runtime.main

# Always come back up
Restart=always
RestartSec=5
# If it crashes five times in 30s, back off — avoids thrash on a bad config
StartLimitIntervalSec=30
StartLimitBurst=5

# Journald captures stdout/stderr automatically
StandardOutput=journal
StandardError=journal

# Lock it down a touch
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/ai-trading-fund/state /opt/ai-trading-fund/logs /opt/ai-trading-fund/vault

[Install]
WantedBy=multi-user.target
